# React_Setup/raspberry_pi_scripts/startup_scripts/startup_service.service
[Unit]
After=network.service

[Service]
ExecStart=/bin/bash /usr/bin/startup_service.sh

[Install]
WantedBy=default.target
